# deployment/server/moj-warsztat-backup-check.service
[Unit]
Description=Mój Warsztat - kontrola repozytorium restic

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/moj-warsztat-backup-check
